# webservice/crontab
# Added by Alex Hancock
*/2 * * * * . /app/log/env.txt python /app/luigi-interface/new_monitor.py >> /app/log/output.txt 2> /app/log/error.txt

